# build/hfs/runit/metadata/run
#!/bin/sh
exec chpst -u app node /app/bin/dist/components/metadata/bin/server \
    --config /app/conf/metadata.hjson
