# build/webhooks/runit/webhooks/run
#!/bin/bash

export NODE_ENV=production
export NODE_PATH=/app/bin/

PRYV_BOILER_SUFFIX="-$$" exec chpst -u app node /app/bin/components/webhooks/bin/server \
    --config /app/conf/webhooks.yml
