# build/webhooks/runit/app/run
#!/bin/sh
exec chpst -u app node /app/bin/dist/components/webhooks/bin/server \
    --config /app/conf/webhooks.json
