# build/webhooks/runit/app/run
#!/bin/sh
PRYV_BOILER_SUFFIX="-$$" exec chpst -u app node /app/bin/components/webhooks/bin/server \
    --config /app/conf/webhooks.yml
