# build/preview/runit/app/run
#!/bin/sh
exec chpst -u app node /app/bin/dist/components/previews-server/src/server \
    --config /app/conf/preview.json
