# build/preview/runit/app/run
#!/bin/sh
PRYV_BOILER_SUFFIX="-$$" exec chpst -u app node /app/bin/dist/components/previews-server/src/server \
    --config /app/conf/preview.yml
