# build/core/runit/app/run
#!/bin/sh
exec chpst -u app node /app/bin/dist/components/api-server/bin/server \
    --config /app/conf/core.yml --http-port ${PORT_NUM}
# ${PORT_NUM} will be replaced by launching script 